# profile.d/dotnet-cli-telemetry-opt-out.sh
# README: Place this file in /etc/profile.d
# https://aka.ms/dotnet-cli-telemetry
export DOTNET_CLI_TELEMETRY_OPTOUT=1
